# systemd-sockets.socket
[Unit]
Description=The listen socket

[Socket]
ListenStream=8000
ListenStream=8001

[Install]
WantedBy=sockets.target
